# deploy/hetzner/systemd/mcp-backup.service
# Daily encrypted backup of mcp-approval2 + mcp-knowledge2 + Vault.
#
# Installed by deploy/hetzner/setup.sh into /etc/systemd/system/.
# Triggered by mcp-backup.timer (03:00 UTC). Logs land in journalctl
# under unit mcp-backup.service.

[Unit]
Description=mcp-approval2 daily encrypted backup
Documentation=file:///opt/mcp-approval2/docs/runbooks/runbook-hetzner-backup-restore.md
# Don't fire while a deploy is in-flight (update.sh holds this same
# state implicitly — pre-deploy backup runs once before update.sh).
ConditionPathExists=/opt/mcp-approval2/deploy/hetzner/.env
ConditionPathExists=/opt/mcp-approval2/.backup-key
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=deploy
Group=deploy
WorkingDirectory=/opt/mcp-approval2/deploy/hetzner
ExecStart=/usr/bin/env bash /opt/mcp-approval2/deploy/hetzner/backup.sh
# Soft-cap a stuck backup so we don't hold the timer hostage. 30 min
# is generous — a full dump of approval2+knowledge2 at pilot size is
# typically under 2 min.
TimeoutStartSec=1800
StandardOutput=journal
StandardError=journal
# Retry once on transient failures (e.g. transient docker hiccup).
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target
